# deploy/docker-test-setup/ebpf-guard-test.service
[Unit]
Description=ebpf-guard — eBPF runtime security agent (VPS test setup)
Documentation=https://github.com/zugolO/ebpf-guard
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/ebpf-guard
ExecStart=/opt/ebpf-guard/build/ebpf-guard --config=/opt/ebpf-guard/deploy/docker-test-setup/config-test.yaml --log-level=info
Restart=on-failure
RestartSec=10
LimitNOFILE=65536
LimitMEMLOCK=infinity

RuntimeDirectory=ebpf-guard
ReadWritePaths=/var/lib/ebpf-guard /var/log/ebpf-guard

[Install]
WantedBy=multi-user.target
